# monitoring.service
[Unit]
Description=Мониторим процесс тест
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh
